# homeassistant-addon/rootfs/etc/services.d/calendifier/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Calendifier service finish script
# ==============================================================================

bashio::log.info "Calendifier service stopped"

# Clean shutdown
exit 0