# homeassistant-addon/rootfs/etc/services.d/calendifier/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start Calendifier service
# ==============================================================================

# Get configuration
CONFIG_PATH="/data/options.json"
LOG_LEVEL=$(bashio::config 'log_level')
LOCALE=$(bashio::config 'locale')
TIMEZONE=$(bashio::config 'timezone')
HOLIDAY_COUNTRY=$(bashio::config 'holiday_country')
THEME=$(bashio::config 'theme')
SSL=$(bashio::config 'ssl')

# Set log level
bashio::log.info "Starting Calendifier with log level: ${LOG_LEVEL}"

# Auto-detect locale if set to auto
if [ "${LOCALE}" = "auto" ] || [ -z "${LOCALE}" ]; then
    LOCALE=$(bashio::info.language)
    bashio::log.info "Auto-detected locale: ${LOCALE}"
fi

# Auto-detect timezone if set to auto
if [ "${TIMEZONE}" = "auto" ] || [ -z "${TIMEZONE}" ]; then
    TIMEZONE=$(bashio::info.timezone)
    bashio::log.info "Auto-detected timezone: ${TIMEZONE}"
fi

# Auto-detect holiday country if set to auto
if [ "${HOLIDAY_COUNTRY}" = "auto" ] || [ -z "${HOLIDAY_COUNTRY}" ]; then
    # Extract country from locale (e.g., en_US -> US, de_DE -> DE)
    HOLIDAY_COUNTRY=$(echo "${LOCALE}" | cut -d'_' -f2)
    bashio::log.info "Auto-detected holiday country: ${HOLIDAY_COUNTRY}"
fi

# Set environment variables
export CALENDIFIER_LOG_LEVEL="${LOG_LEVEL}"
export CALENDIFIER_LOCALE="${LOCALE}"
export CALENDIFIER_TIMEZONE="${TIMEZONE}"
export CALENDIFIER_HOLIDAY_COUNTRY="${HOLIDAY_COUNTRY}"
export CALENDIFIER_THEME="${THEME}"
export CALENDIFIER_DATA_DIR="/data/calendifier"
export CALENDIFIER_CONFIG_FILE="/data/calendifier/settings.json"
export CALENDIFIER_SSL="${SSL}"

# Create data directory
mkdir -p /data/calendifier/data
mkdir -p /data/calendifier/logs
mkdir -p /data/calendifier/exports
mkdir -p /data/calendifier/backups

# Set permissions
chown -R root:root /data/calendifier

# Start Calendifier API server
bashio::log.info "Starting Calendifier API server..."
cd /app

# Use SSL if enabled
if bashio::config.true 'ssl'; then
    CERTFILE=$(bashio::config 'certfile')
    KEYFILE=$(bashio::config 'keyfile')
    
    bashio::log.info "SSL enabled, using certificates: ${CERTFILE}, ${KEYFILE}"
    
    exec python3 -m calendifier_addon.api_server \
        --host 0.0.0.0 \
        --port 8099 \
        --ssl-certfile "/ssl/${CERTFILE}" \
        --ssl-keyfile "/ssl/${KEYFILE}" \
        --log-level "${LOG_LEVEL}"
else
    bashio::log.info "Starting without SSL"
    
    exec python3 -m calendifier_addon.api_server \
        --host 0.0.0.0 \
        --port 8099 \
        --log-level "${LOG_LEVEL}"
fi